# docker/ferrum.service
[Unit]
Description=Ferrum WhatsApp Orchestrator
After=network.target

[Service]
Type=simple
User=ferrum
WorkingDirectory=/opt/ferrum
EnvironmentFile=/opt/ferrum/.env
ExecStart=/opt/ferrum/ferrum run
Restart=always
RestartSec=3
LimitNOFILE=65536
NoNewPrivileges=true
ProtectSystem=full
ProtectHome=true
PrivateTmp=true

[Install]
WantedBy=multi-user.target
